# mkosi.presets/initrd/mkosi.extra/usr/lib/systemd/system/installer.target
[Unit]
Description=Installer Mode
Requires=sysinit.target installer.service
After=sysinit.target installer.service
AllowIsolate=yes
